# saltstack/base/salt/app/api/files/supervisor.ini
[program:api-webapp]
command=/usr/run/vk-fetch/api/current/pyenv.sh gunicorn settings.wsgi -c /usr/run/vk-fetch/api/current/settings/gunicorn.conf.py
user=user
environment=HOME="/home/user",USER="user",DJANGO_SETTINGS_MODULE="{{ pillar['django_conf'][grains['environment']] }}",ENVIRONMENT="{{ grains['environment'] }}"
directory=/usr/run/vk-fetch/api/current
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/api-webapp.log

[program:celery-beat]
command=/usr/run/vk-fetch/api/current/pyenv.sh ./manage.py celery beat -l info --pidfile=/usr/run/vk-fetch/celerybeat.pid
user=user
environment=HOME="/home/user",USER="user",DJANGO_SETTINGS_MODULE="{{ pillar['django_conf'][grains['environment']] }}",ENVIRONMENT="{{ grains['environment'] }}"
directory=/usr/run/vk-fetch/api/current
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/celery-beat.log


[program:celery-workers]
command=/usr/run/vk-fetch/api/current/pyenv.sh ./manage.py celery worker -l info -c 4
user=user
environment=HOME="/home/user",USER="user",DJANGO_SETTINGS_MODULE="{{ pillar['django_conf'][grains['environment']] }}",ENVIRONMENT="{{ grains['environment'] }}"
directory=/usr/run/vk-fetch/api/current
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/celery-worker.log
