# saltstack/base/salt/app/api/files/supervisor.ini
[program:api]
command=/home/user/concierge/concierge-code/current/src/pyenv gunicorn settings.wsgi -c /home/user/concierge/concierge-code/current/src/settings/gunicorn.conf.py
user=user
environment=HOME="/home/user",USER="user"
directory=/home/user/concierge/concierge-code/current/src
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/user/concierge/logs/agent-webapp.log
